# fls_ws.service
[Unit]
Description=FeliCa WebSocket Server
After=network.target pcscd.service
Wants=pcscd.service

[Service]
Type=simple
# 実行ファイルへのパス（必要に応じて -p オプション等を指定してください）
ExecStart=/usr/local/bin/fls_ws -p
WorkingDirectory=/usr/local/bin

# サービスを実行するユーザー（PC/SCリーダーへのアクセス権限を持つユーザーを指定）
# User=igari

# クラッシュ時の自動再起動設定
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
